# angstromctf/2020/web/defunds_crypt/cleardir.cron
SHELL=/bin/bash
30 * * * * rm /var/www/memories/*

